# .config/systemd/user/rclone.service
[Unit]
Description=Use RClone to Backup Files

[Service]
Type=simple
ExecStart=rclone -v sync ~/org Google:/org
